# docker/prod/queueworker/conf/supervisord/callCredit-queue.conf
[program:callCredit]
command=/bin/bash /var/www/html/callCredit-queue.sh
startsecs=0
autorestart=false
startretries=0

