# docker/prod/queueworker/conf/supervisord/propertyLive-queue.conf
[program:propertyLive]
command=/bin/bash /var/www/html/propertyLive-queue.sh
startsecs=0
autorestart=false
startretries=0

