# docker/prod/queueworker/conf/supervisord/sun-queue.conf
[program:sun]
command=/bin/bash /var/www/html/sun-queue.sh
startsecs=3
autorestart=true
startretries=30
