# docker/prod/queueworker/conf/supervisord/trustId-queue.conf
[program:trustId]
command=/bin/bash /var/www/html/trustId-queue.sh
startsecs=0
autorestart=false
startretries=0

