# docker/dev/queueworker/conf/supervisord/callCredit-queue.conf
[program:callCredit]
command=/bin/bash /var/www/html/callCredit-queue.sh
startsecs=0
autorestart=true
startretries=0
