# docker/dev/queueworker/conf/supervisord/docusign.conf
[program:docusign]
command=/bin/bash /var/www/html/docusign.sh
startsecs=0
autorestart=true
startretries=0
