# docker/prod/portal/conf/supervisord/filebeat.conf
[program:filebeat]
command=filebeat -e -c /etc/filebeat/conf.d/filebeat.yml  -path.data /etc/filebeat/conf.d/data
redirect_stderr=true
